# com.blw.mqtt-logger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blw.mqtt-logger</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/Files/claude/mqtt-logger/venv/bin/python3</string>
        <string>/Volumes/Files/claude/mqtt-logger/mqtt_logger.py</string>
        <string>--broker</string>
        <string>localhost</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/Files/claude/mqtt-logger</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
